# bot/config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:xvfb]
command=Xvfb :99 -screen 0 1280x720x24 -ac +extension GLX +render -noreset
autorestart=true
priority=100
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb-err.log
environment=HOME="/root"

[program:dbus]
command=dbus-daemon --system --nofork --nopidfile
autorestart=true
priority=150
stdout_logfile=/var/log/supervisor/dbus.log
stderr_logfile=/var/log/supervisor/dbus-err.log
environment=HOME="/root"
startsecs=1

[program:pulseaudio]
command=pulseaudio --system --daemonize=no --exit-idle-time=-1 --disallow-exit --log-level=info --log-target=stderr --file=/etc/pulse/system.pa
autorestart=true
priority=200
stdout_logfile=/var/log/supervisor/pulseaudio.log
stderr_logfile=/var/log/supervisor/pulseaudio-err.log
environment=HOME="/root",DISPLAY=":99",XDG_RUNTIME_DIR="/tmp"
startsecs=2

[program:bot]
command=python3 /app/src/bot.py
autorestart=false
priority=300
stdout_logfile=/var/log/supervisor/bot.log
stderr_logfile=/var/log/supervisor/bot-err.log
environment=HOME="/root",DISPLAY=":99",PULSE_SERVER="unix:/tmp/pulse/native",XDG_RUNTIME_DIR="/tmp"
startsecs=5